# crontab
# START CRON JOB
*/20 * * * * /usr/local/bin/python /app/download_prices.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
